# systemd/eyebrod.service
[Unit]
Description=eyebrod

[Service]
ExecStart=/usr/bin/eyebrod
Restart=always

[Install]
WantedBy=basic.target
